# scripts/update-gcloud-image/update-gcloud-image.timer
[Unit]
Description=Run update-gcloud-image once a day

[Timer]
OnBootSec=15min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
